# yale-doorman.service
[Unit]
Description=Yale Doorman L3S BLE Monitor
After=bluetooth.target network.target
Wants=bluetooth.target

[Service]
Type=simple
User=alex
WorkingDirectory=/home/alex/vscode
ExecStart=/home/alex/vscode/venv/bin/python3 /home/alex/vscode/app.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="PATH=/home/alex/vscode/venv/bin:/usr/bin:/bin"

# Security hardening
NoNewPrivileges=true
ProtectHome=read-only
ReadWritePaths=/home/alex/.yale-doorman

[Install]
WantedBy=multi-user.target
